# services/mbsync.timer
[Unit]
Description=Sync mail (run each 5 minutes)

[Timer]
OnBoot=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
