# services/mbsync.service
[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -a
ExecStart=notmuch new
