# services/.config/systemd/user/mbsync.timer
[Unit]
Description=Sync mail (run daily)

[Timer]
#OnCalendar=*:0/10
OnCalendar=*-*-* 00:00:00


[Install]
WantedBy=timers.target
