# services/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=Sync vdir (Calendars and contacts) (run each 30 minutes)

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
