# services/.config/systemd/user/vdirsyncer.service
[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
